# deploy/auto-agent.service
[Unit]
Description=Auto Agent Harness UI Server
Documentation=https://github.com/nek1987/auto-agent-harness
After=network.target

[Service]
Type=simple
User=<username>
Group=<group>
WorkingDirectory=/path/to/auto-agent-harness

# Environment
Environment="PATH=/path/to/auto-agent-harness/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="HOST=0.0.0.0"
Environment="PORT=8888"

# Main command
ExecStart=/path/to/auto-agent-harness/venv/bin/python -m uvicorn server.main:app --host 0.0.0.0 --port 8888

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Logging
StandardOutput=append:/var/log/auto-agent/server.log
StandardError=append:/var/log/auto-agent/error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
